# application.ini
[program:fplpartner]
command=/source/cms-comment
directory=/source
redirect_stderr=true
stdout_logfile=/data/log/fplpartner.log
stdout_logfile_maxbytes=0
